# rootfs.extra/etc/inittab
# Based on the default busybox config
::sysinit:/etc/init.d/rcS
::askfirst:-/bin/sh
::ctrlaltdel:/sbin/reboot
::restart:/sbin/init
tty1::askfirst:-/bin/sh
